# etc/systemd/system/bt-kbd.service
[Service]
Type=oneshot
User=mag
WorkingDirectory=/home/mag
ExecStart=/home/mag/bin/bt-kbd.sh
